# crontab
# Cron configuration for snapshots every minute
* * * * * /app/snapshot >> /var/log/cron.log 2>&1

# Empty line required at the end of cron file